# base/etc/init/auto-frontend.conf
start on desktop-session-start or started vdr or started openbox or resume

task

script
  . /usr/share/yavdr/helpers/upstart-base
  start_frontend
end script
