# softhddevice/etc/init/softhddevice-suspend.conf
start on yavdr-stop-frontend

task

exec vdr-dbus-send /Plugins/softhddevice plugin.SVDRPCommand string:'SUSP'
